# deploy/pi/systemd/f2a-web-ui.service
[Unit]
Description=first2apply web app (Next.js standalone, PWA)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=first2apply
Group=first2apply
WorkingDirectory=/opt/first2apply/builds/web-ui
EnvironmentFile=/opt/first2apply/.env
Environment=PORT=3030
# Bind on all interfaces so the dockerized maadcloud-caddy container can
# reach the service via the docker bridge gateway (172.17.0.1). Loopback-only
# binding would be invisible to containers. The Pi sits behind a NAT and only
# Caddy holds :443; raw :3030 is not exposed off-host.
Environment=HOSTNAME=0.0.0.0
# Next.js standalone server entry. The deploy script rsyncs
# apps/webapp/.next/standalone/ into WorkingDirectory, so this resolves to
# /opt/first2apply/builds/web-ui/apps/webapp/server.js.
ExecStart=/usr/bin/node apps/webapp/server.js
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/first2apply
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
