# deploy/pi/systemd/f2a-web-ui.service
[Unit]
Description=first2apply web app (Next.js standalone, PWA)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=first2apply
Group=first2apply
WorkingDirectory=/opt/first2apply/builds/web-ui
EnvironmentFile=/opt/first2apply/.env
Environment=PORT=3030
Environment=HOSTNAME=127.0.0.1
# Next.js standalone server entry. The deploy script rsyncs
# apps/webapp/.next/standalone/ into WorkingDirectory, so this resolves to
# /opt/first2apply/builds/web-ui/apps/webapp/server.js.
ExecStart=/usr/bin/node apps/webapp/server.js
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/first2apply
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
