# deploy/pi/systemd/f2a-web-ui.service
[Unit]
Description=first2apply server web UI (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=first2apply
Group=first2apply
WorkingDirectory=/opt/first2apply/builds/web-ui
EnvironmentFile=/opt/first2apply/.env
Environment=PORT=3030
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start -p ${PORT}
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/first2apply
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
